# ops/launchd/com.capas.sota-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.capas.sota-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/kreniq/Library/Application Support/CAPAS/capas-inteligentes' &amp;&amp; /usr/bin/git fetch origin main &amp;&amp; /usr/bin/git reset --hard origin/main &amp;&amp; echo "[$(/bin/date -u +%Y-%m-%dT%H:%M:%SZ)] CAPAS daily SOTA watch started" &amp;&amp; /usr/bin/python3 scripts/update_sota_watch.py &amp;&amp; echo "[$(/bin/date -u +%Y-%m-%dT%H:%M:%SZ)] CAPAS validation started" &amp;&amp; /usr/bin/python3 capas.py validate &amp;&amp; echo "[$(/bin/date -u +%Y-%m-%dT%H:%M:%SZ)] CAPAS daily SOTA watch finished"</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kreniq/Library/Application Support/CAPAS/capas-inteligentes</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/kreniq/Library/Application Support/CAPAS/capas-inteligentes/outputs/logs/sota-daily.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kreniq/Library/Application Support/CAPAS/capas-inteligentes/outputs/logs/sota-daily.err.log</string>
</dict>
</plist>
